# config/includes.chroot/etc/systemd/system/kiosk-selftest.service
[Unit]
Description=One-shot kiosk boot self-test
After=kiosk.service
ConditionKernelCommandLine=kiosk.selftest=1

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/kiosk-selftest

[Install]
WantedBy=graphical.target
